# scripts/liotad.socket
[Socket]
ListenFIFO=/var/tmp/liota/package_messenger.fifo

[Install]
WantedBy=sockets.target